# mail.cron
#!/bin/bash
notmuch new| /usr/bin/logger
/home/pietrushnic/bin/afew -t -vvv -n | /usr/bin/logger
